# snaptransfer/snapxfer/crontab/snapxfer.crontab
0,10,20,30,40,50 * * * * /opt/custom/snapxfer/snapxfer async
3 * * * * /opt/custom/snapxfer/snapxfer hourly
7 1 * * * /opt/custom/snapxfer/snapxfer daily
13 3 * * 1 /opt/custom/snapxfer/snapxfer weekly
17 7 5 * * /opt/custom/snapxfer/snapxfer monthly
