# crontab/crontab/zfsscrub.crontab
1 1 1 * * /opt/custom/zfsscrub/scrub_and_report.sh
